# pkg/ghoneycutt-logrotate-1.0.0/files/logrotate.cron
#!/bin/sh
# This file is being maintained by Puppet.
# DO NOT EDIT

/usr/sbin/logrotate /etc/logrotate.conf
EXITVALUE=$?
if [ $EXITVALUE != 0 ]; then
    /usr/bin/logger -t logrotate "ALERT exited abnormally with [$EXITVALUE]"
fi
exit 0
